# starlink_remoto/systemd/starlink-rele@.service
[Unit]
Description=Control del rele de Starlink (%i)

[Service]
Type=oneshot
ExecStart=/root/starlink_remoto/control_starlink.sh %i
